# systemd/comfyui.service
[Unit]
Description=ComfyUI service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/ComfyUI
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/ubuntu/ComfyUI/comfy-env/bin/python /home/ubuntu/ComfyUI/main.py --listen 127.0.0.1 --port 8188
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target